# deploy/labgate-sweep.cron.example
# Install in root's crontab only after creating root-owned mode-0600
# /etc/labgate/cron-curl.conf as documented in docs/OPERATIONS.md. Keeping the
# bearer value in that file prevents it from appearing in process listings.
* * * * * curl --config /etc/labgate/cron-curl.conf --fail --silent --show-error --max-time 20 --request POST http://127.0.0.1:3000/api/cron/sweep >/dev/null
