# machine-setup/guest-cleanup.service
[Unit]
Description=Expire an abandoned LabGate guest session
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/guest-cleanup.sh
