# machine-setup/guest-webhook-flush.path
[Unit]
Description=Flush queued LabGate lifecycle webhooks when the outbox changes

[Path]
PathChanged=/var/lib/labgate/outbox
Unit=guest-webhook-flush.service

[Install]
WantedBy=multi-user.target
