# machine-setup/guest-cleanup.service
[Unit]
Description=Expire or recover a LabGate guest credential
After=local-fs.target
RequiresMountsFor=/home/guest /var/lib/labgate

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/guest-cleanup.sh
UMask=0077
TimeoutStartSec=30s
