# machine-setup/guest-heartbeat.service
[Unit]
Description=Send a LabGate machine heartbeat
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/guest-heartbeat.sh
