# machine-setup/guest-cleanup.timer
[Unit]
Description=Check for abandoned LabGate guest sessions

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=5s
Persistent=true
Unit=guest-cleanup.service

[Install]
WantedBy=timers.target
